# init/systemd/buteo-mtp.service
[Unit]
Description=MTP service

[Service]
ExecStart=/usr/bin/mtpd
Restart=always
RestartSec=1
Nice=-10
